# deploy/toloseo.service
[Unit]
Description=Toloseo Backend
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/toloseo/server
# --smol halves Bun's heap baseline: safer on a 2G VPS next to GTFS in-memory store.
ExecStart=/usr/local/bin/bun --smol run src/index.ts
Restart=always
RestartSec=5
# GTFS static parse can take 60-90s on cold cache; give bootstrap headroom.
TimeoutStartSec=120
Environment=NODE_ENV=production
Environment=PORT=3001
StandardOutput=journal
StandardError=journal
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
