# deploy/toloseo.service
[Unit]
Description=Toloseo Backend
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/toloseo/server
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001
StandardOutput=journal
StandardError=journal
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
